# conf/systemd.service
[Unit]
Description=__APP__ Go Application
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__/
Environment="PORT=__PORT__"
Environment="DOMAIN=__DOMAIN__"
Environment="PATH_URL=__PATH__"
Environment="WEBHOOK_KEY=__WEBHOOK_KEY__"
Environment="SENDMAIL_PATH=/usr/sbin/sendmail"
Environment="BACKEND_TYPE=__BACKEND_TYPE__"
Environment="SMTP_HOST=__SMTP_HOST__"
Environment="SMTP_PORT=__SMTP_PORT__"
Environment="SMTP_USER=__SMTP_USER__"
Environment="SMTP_PASS=__SMTP_PASS__"
ExecStart=__INSTALL_DIR__/bin/__APP__
StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=inherit

# Sandboxing options - ULTRA RELAXED FOR DEBUGGING
# NoNewPrivileges must be false for sendmail/postdrop to use setgid
NoNewPrivileges=false
PrivateTmp=false
PrivateDevices=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=false
RestrictRealtime=true
DevicePolicy=open
ProtectSystem=false
ProtectHome=false
ReadWritePaths=__DATA_DIR__ /var/spool/postfix

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
