# conf/systemd.service
[Unit]
Description=__APP__ Go Application
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__/
Environment="PORT=__PORT__"
Environment="DOMAIN=__DOMAIN__"
Environment="PATH_URL=__PATH__"
Environment="WEBHOOK_KEY=__WEBHOOK_KEY__"
Environment="SENDMAIL_PATH=/usr/sbin/sendmail"
ExecStart=__INSTALL_DIR__/bin/__APP__
StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=inherit

# Sandboxing options
# NoNewPrivileges must be false for sendmail/postdrop to use setgid
NoNewPrivileges=false
PrivateTmp=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=false
RestrictRealtime=true
DevicePolicy=closed
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__DATA_DIR__

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
